# root/etc/services.d/app/run
#!/usr/bin/with-contenv sh
whoami

chmod +x /app/complier

cd /app && exec /app/complier;